# deploy/datalake-refresh.timer
[Unit]
Description=Trigger datalake weekly MT5 refresh

[Timer]
# Sundays at 03:00 UTC — runs before the backup timer so the weekend's fresh data gets exported.
OnCalendar=Sun *-*-* 03:00:00 UTC
Persistent=true
RandomizedDelaySec=5min
Unit=datalake-refresh.service

[Install]
WantedBy=timers.target
